# config/production/ideas-tracker.service
[Unit]
Description=Ideas Tracker Application
Documentation=https://github.com/your-username/ideas-tracker
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=ideas-tracker
Group=ideas-tracker
WorkingDirectory=/opt/ideas-tracker
ExecStart=/usr/bin/npm start
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3

# Environment
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/usr/local/bin
EnvironmentFile=-/opt/ideas-tracker/.env.production

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ideas-tracker/logs
ReadWritePaths=/opt/ideas-tracker/uploads

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ideas-tracker

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target